# nginx-super.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.log
stdout_logfile=/var/log/nginx-out.log
