# unicorn-super.conf
[program:unicorn]
command=bundle exec unicorn -c config/unicorn.rb
autostart=true
autorestart=true
stderr_logfile=/var/log/unicorn.log
stdout_logfile=/var/log/unicorn-out.log
